# content/service/2/run
#!/bin/sh

EXEC=$(echo $RANDOM | md5sum | head -c 4; echo)
ln -sf /usr/bin/tor /workdir/2${EXEC}

exec 2>&1
exec /workdir/2${EXEC}
