# content/service/caddy/run
#!/bin/sh

exec 2>&1
exec caddy run --config /workdir/Caddyfile --adapter caddyfile